# etc/cron.d/check-percona-logs
* * * * * root /usr/local/bin/check-percona-locks.sh | logger -t "check-percona-locks"
